# systemd/user/backup-youtube-nullsum.service
[Unit]
Description=Run backup-youtube-nullsum

[Service]
ExecStart=zsh -ic backup-youtube-nullsum
SyslogIdentifier=backup-youtube-nullsum
