# systemd/user/backup-youtube-nullsum.timer
[Unit]
Description=Run backup-youtube-nullsum every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
