# systemd/user/certbot-renew.service
[Unit]
Description=Certbot renew

[Service]
ExecStart=/home/ar/src/certs-nullsum.net/renew.sh
SyslogIdentifier=certbot-renew
