# systemd/user/certbot-renew.service
[Unit]
Description=Certbot renew

[Service]
ExecStart=zsh -ic /home/ar/Code/certs-nullsum.net/renew.sh
SyslogIdentifier=certbot-renew
