# systemd/user/check-update-phobos-bios.service
[Unit]
Description=Run check-update-phobos-bios

[Service]
ExecStart=zsh -lc check-update-phobos-bios
SyslogIdentifier=check-update-phobos-bios
