# systemd/user/check-update-mars-bios.service
[Unit]
Description=Run check-update-mars-bios

[Service]
ExecStart=zsh -lc check-update-mars-bios
SyslogIdentifier=check-update-mars-bios
