# systemd/user/download-latest-vuescan-version.timer
[Unit]
Description=Run download-latest-vuescan-version every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
