# systemd/user/git-repo-backup.service
[Unit]
Description=Run git-repo-backup

[Service]
ExecStart=zsh -lc git-host-backup
SyslogIdentifier=git-repo-backup
