# systemd/user/git-repo-backup.service
[Unit]
Description=Run git-repo-backup

[Service]
ExecStart=/usr/bin/env backup-git-repos
SyslogIdentifier=git-repo-backup
