# systemd/user/github-repo-backup.service
[Unit]
Description=Run github-repo-backup

[Service]
ExecStart=/usr/bin/env github-repo-backup
SyslogIdentifier=github-repo-backup
