# systemd/user/gog-update.service
[Unit]
Description=Run gog-update

[Service]
ExecStart=/usr/bin/env gog-update
SyslogIdentifier=gog-update
