# systemd/user/gog-update.service
[Unit]
Description=Run gog-update

[Service]
ExecStart=zsh -lc gog-update
SyslogIdentifier=gog-update
