# systemd/user/gog-update.service
[Unit]
Description=Run gog-update

[Service]
ExecStart=/home/ar/.dotfiles/public/scripts/hosts/sol/gog-update
SyslogIdentifier=gog-update
