# systemd/user/gog-update.timer
[Unit]
Description=Run gog-update every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
