# systemd/user/harmonize-library.service
[Unit]
Description=Run harmonize-library

[Service]
ExecStart=harmonize-library
SyslogIdentifier=harmonize-library
