# systemd/user/harmonize-library.timer
[Unit]
Description=Run harmonize-library every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
