# systemd/user/harmonize-library.service
[Unit]
Description=Run harmonize-library

[Service]
ExecStart=zsh -ic harmonize-library
SyslogIdentifier=harmonize-library
