# systemd/user/linode-dynamic-dns.service
[Unit]
Description=Run linode-dynamic-dns

[Service]
ExecStart=/usr/bin/linode-dynamic-dns
