# systemd/user/mbsync.service
[Unit]
Description=Run mbsync

[Service]
ExecStart=zsh -lc 'mbsync -a'
SyslogIdentifier=mbsync
