# systemd/user/mbsync.service
[Unit]
Description=Run mbsync

[Service]
ExecStart=/usr/bin/env mbsync -a
SyslogIdentifier=mbsync
