# systemd/user/mbsync.service
[Unit]
Description=Run mbsync

[Service]
ExecStart=mbsync -a
SyslogIdentifier=mbsync
