# systemd/user/mpd-notifier.service
[Unit]
Description=mpd-notifier

[Service]
ExecStart=zsh -lc 'exec mpd-notifier'
SyslogIdentifier=mpd-notifier
Restart=on-failure
RestartSec=5s
