# systemd/user/retroarch-rom-sync.timer
[Unit]
Description=Run retroarch-rom-sync every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
