# systemd/user/retroarch-rom-sync.service
[Unit]
Description=Run retroarch-rom-sync

[Service]
ExecStart=zsh -ic retroarch-rom-sync
SyslogIdentifier=retroarch-rom-sync
