# systemd/user/sol-update-remotes.service
[Unit]
Description=Run sol-update-remotes

[Service]
ExecStart=sol-update-remotes
SyslogIdentifier=sol-update-remotes
