# systemd/user/sol-ytdl-sub.service
[Unit]
Description=Run sol-ytdl-sub

[Service]
ExecStart=zsh -lc sol-ytdl-sub
SyslogIdentifier=sol-ytdl-sub
