# systemd/user/sol-sync-camera3.timer
[Unit]
Description=Run sol-sync-camera3

[Timer]
OnCalendar=*-*-* 4:00:00
Persistent=true

[Install]
WantedBy=timers.target
