# systemd/user/sol-fix-permissions.service
[Unit]
Description=Run sol-fix-permissions

[Service]
ExecStart=zsh -lc sol-fix-permissions
SyslogIdentifier=sol-fix-permissions
