# systemd/user/sol-update-pacman-mirrorlist.timer
[Unit]
Description=Run sol-update-pacman-mirrorlist

[Timer]
OnCalendar=Wed *-* 5:00:00
Persistent=true

[Install]
WantedBy=timers.target
