# systemd/user/sol-fix-permissions.timer
[Unit]
Description=Run sol-fix-permissions every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
