# systemd/user/sol-status-email.service
[Unit]
Description=Run sol-status-email

[Service]
ExecStart=zsh -lc sol-status-email
EnvironmentFile=/home/ar/.dotfiles/private/sol-status-email.env
SyslogIdentifier=sol-status-email
