# systemd/user/sol-update-pacman-mirrorlist.service
[Unit]
Description=Run sol-update-pacman-mirrorlist

[Service]
ExecStart=zsh -ic pacman-update-mirrors
SyslogIdentifier=sol-update-pacman-mirrorlist
