# systemd/user/sol-ytdl-sub.timer
[Unit]
Description=Run sol-ytdl-sub every 30 minutes

[Timer]
OnBootSec=30min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
