# systemd/user/sol-update-check.service
[Unit]
Description=Run sol-update-check

[Service]
ExecStart=zsh -lc sol-update-check
SyslogIdentifier=sol-update-check
