# systemd/user/sol-status-email.service
[Unit]
Description=Run sol-status-email

[Service]
ExecStart=zsh -lc sol-status-email
SyslogIdentifier=sol-status-email
