# systemd/user/sol-sync-camera1.service
[Unit]
Description=Run sol-sync-camera1

[Service]
ExecStart=zsh -lc sync-camera1.nullsum.lan --format-storage
SyslogIdentifier=sol-sync-camera1
