# systemd/user/wasabi-sync.service
[Unit]
Description=Run wasabi-sync

[Service]
ExecStart=wasabi-sync
SyslogIdentifier=wasabi-sync
