# systemd/user/wasabi-sync.timer
[Unit]
Description=Run wasabi-sync every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
