# non-user/sol/etc/systemd/system/restic-backup.service
[Unit]
Description=Run restic-backup

[Service]
SyslogIdentifier=restic-backup
ExecStart=/home/ar/.dotfiles/private/scripts/terminal/restic-backup
