# non-user/sol/etc/systemd/system/sol-backup.service
[Unit]
Description=Run sol-backup

[Service]
SyslogIdentifier=sol-backup
ExecStart=/home/ar/.dotfiles/private/scripts/hosts/sol/sol-backup
