# non-user/sol/etc/systemd/system/sol-backup.timer
[Unit]
Description=Run sol-backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
